# release/installer/deb/postinst
#!/bin/sh
# postinst script for pivotal-gemfirexd
#
# see: dh_installdeb(1)

set -e

case "$1" in
  configure)

    /bin/mkdir -p /var/log/pivotal/gemfirexd
    chown -R gemfirexd:pivotal /var/log/pivotal/gemfirexd

    # Add the init script to chkconfig/update-rc.d
    if [ -x /sbin/chkconfig ]; then
      /sbin/chkconfig --add gemfirexd
    elif [ -x /usr/sbin/update-rc.d ]; then
      /usr/sbin/update-rc.d gemfirexd defaults 85 15
    fi
    
    # Add symlinks from usr/bin to executable scripts
    ln -sf /opt/pivotal/gemfirexd/@INSTALL_DIR@/bin/gfxd /usr/bin/gfxd

    ;;

    *)
        #echo "postinst called with unknown argument \`$1'" >&2
        #exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
